# etc/systemd/system/user@60466.service.d/delegate.conf
[Service]
Delegate=cpuset cpu io memory pids